# crontab
* * * * *  . /etc/profile; python3 CrontabController.py